# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -d
startsecs = 0
autorestart = true

[program:django]
command=/usr/local/bin/python /src/manage.py runserver 0.0.0.0:8000

